# nginx-php-exim/etc/service/nginx/run
#!/bin/bash

# exit on error
set -e

# Start nginx
exec nginx
