# ubuntu/etc/rc.local
#!/usr/bin/env bash
set -e

[[ ! -d /etc/rc.docker.d ]] && mkdir -p /etc/rc.docker.d

# Configure environment
shopt -s nullglob dotglob
for rc in /etc/rc.docker.d/*
do
  # shellcheck source=/dev/null
  echo "Including: $rc"
  . "${rc}"
done
shopt -u nullglob dotglob
