# nginx-hhvm-exim/etc/service/hhvm/run
#!/bin/sh
set -e

exec /usr/bin/hhvm --config /etc/hhvm/server.ini -m server -vPidFile=/var/run/hhvm/pid -vServer.Type=fastcgi
